# supervisord.ini
[supervisord]
nodaemon=true

[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5901 --web /opt/noVNC --file-only
autorestart=true

[program:vnc]
command=vncserver :1
autorestart=true
