# supervisord.ini
[supervisord]
nodaemon=true

[program:TigerVNC]
command=/usr/bin/vncserver :1 -localhost no 
autorestart=true

[program:xfce]
command=/usr/bin/startxfce4
autorestart=false
