# supervisord.ini
[supervisord]
nodaemon=true

[program:x11vnc]
command=vncserver :0
autorestart=true

[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 
autorestart=true

[program:ngrok]
command=ngrok http 6080
autorestart=true
