# supervisord.ini
[supervisord]
nodaemon=true

[program:X11]
command=/usr/bin/Xvfb :0 -screen 0 1024x768x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc
autorestart=true

[program:xfce]
command=/bin/bash -c "startxfce4 && tail -f /dev/null"
autorestart=true
[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900
autorestart=true
