# supervisord.ini
[supervisord]
nodaemon=true

[program:x11vnc]
command=vncserver :0
autorestart=true

[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 
autorestart=true

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork
autorestart=true

[program:xfsettingsd]
command=/usr/bin/xfsettingsd
autorestart=true
