# supervisord.ini
[supervisord]
nodaemon=false

[program:x11vnc]
command=nohup vncserver :0 &
autorestart=false

[program:noVNC]
command=nohup /opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 0.0.0.0:8083 &
autorestart=false

[program:VNC]
command=nohup /root/vnc.sh & 
autorestart=false
