# supervisord.ini
[supervisord]
nodaemon=true

[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --web /opt/noVNC --file-only
autorestart=true

