# supervisord.ini
[supervisord]
nodaemon=true

[program:TigerVNC]
command=/usr/bin/vncserver :0
autorestart=true

[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 0.0.0.0:8083
autorestart=true
