# supervisord.ini
[supervisord]
nodaemon=true

[program:x11vnc]
command=vncserver :0
autorestart=true

[program:noVNC]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 0.0.0.0:8083 
autorestart=true

[program:ngrok]
command=ngrok http --domain=dogfish-intimate-noticeably.ngrok-free.app 8083
autorestart=true
